# files/etc/uci-defaults/09-btrfs-partition
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.
# Uncomment lines to apply:

# additional storage for cosmos? (mmcblk0=internal eMMC)
storage="mmcblk0"

#################

# log potential errors
exec >/tmp/09-btrfs-partition.log 2>&1

#mount eMMC (btrfs) filesystem (used for cosmos)
if [ -n "$storage" ]; then
  echo "add 3rd partition (btrfs) on empty space..."
  # add 3rd partition after 1024MB
  echo "n
  p
  3
  2097153

  w" | fdisk -W never /dev/"$storage"
  #format partition
  echo "format 3rd partition (btrfs)..."
  mkfs.btrfs /dev/"$storage"p3 -f
  # mount btrfs partition to docker dir
  echo "mount 3rd partition (btrfs)..."
  block detect | uci import fstab
  uci set fstab.@mount[-1].target='/opt/docker'
  uci set fstab.@mount[-1].enabled='1'
  uci set fstab.@mount[-1].options='compress=zstd:10'
  uci commit fstab
  #set btrfs docker storage driver
  echo "switch docker storage driver to btrfs..."
  uci set dockerd.globals.storage_driver="btrfs"
  uci commit dockerd
fi

echo "All done!"
exit 0