# files/etc/uci-defaults/09-btrfs-partition
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.

# log potential errors
exec >/opt/09-btrfs-partition.log 2>&1

# additional storage for cosmos? (mmcblk0=internal eMMC)
storage="mmcblk0"

echo "add 3rd partition (btrfs) on empty space"
# add 3rd partition after 1024MB
echo '2097153,,83,' | sfdisk -a -W always --no-reread /dev/"$storage"
echo "partx -u /dev/$storage"
partx -u /dev/"$storage"
echo "partx -u /dev/'$storage'p3"
partx -u /dev/"$storage"p3
#format partition
echo "mkfs.btrfs /dev/'$storage'p3 -f"
mkfs.btrfs /dev/"$storage"p3 -f
#set btrfs docker storage driver
echo "uci set dockerd.globals.storage_driver='btrfs'"
uci set dockerd.globals.storage_driver="btrfs"
# mount btrfs partition to docker dir
echo "block detect | uci import fstab"
block detect | uci import fstab
echo "uci set fstab.@mount[-1].target='/opt/docker'"
uci set fstab.@mount[-1].target='/opt/docker'
echo "uci set fstab.@mount[-1].enabled='1'"
uci set fstab.@mount[-1].enabled='1'
echo "uci set fstab.@mount[-1].options='compress=zstd:10'"
uci set fstab.@mount[-1].options='compress=zstd:10'
echo "uci commit fstab"
uci commit fstab
echo "uci commit dockerd"
uci commit dockerd

echo "All done!"
exit 0